# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[program:nodejs]
command=node /app/server.js
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/nodejs/stderr.log
stdout_logfile=/var/log/nodejs/stdout.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
environment=PORT="9527",DB_HOST="mysql",DB_PORT="3306",DB_USER="root",DB_PASSWORD="root",DB_NAME="deadlock_monitor"

[program:monitor-manager]
command=node /app/monitor-manager.js
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/monitor-manager/stderr.log
stdout_logfile=/var/log/monitor-manager/stdout.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
environment=DB_HOST="mysql",DB_PORT="3306",DB_USER="root",DB_PASSWORD="root",DB_NAME="deadlock_monitor"

[program:namespace-backfill]
command=node /app/namespace-backfill.js
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/namespace-backfill/stderr.log
stdout_logfile=/var/log/namespace-backfill/stdout.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
environment=DB_HOST="mysql",DB_PORT="3306",DB_USER="root",DB_PASSWORD="root",DB_NAME="deadlock_monitor"
